# dockerized/explorer/supervisord-api.conf
[supervisord]
nodaemon=true

[program:api]
command=/usr/local/bin/forever start /rsk-explorer-api/dist/api
autostart=true

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
